# etc/cron.d/mongodb_backup
0 1 * * * root /bin/bash /path/to/gamevault/backend/scripts/backup.sh >> /var/log 2>&1